# solo-container-and-classics-docker/approach1-systemd/client.service
[Unit]
Description=HTTP Client Service
After=network.target server.service
Requires=server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/services
Environment="SERVER_HOST=localhost"
Environment="SERVER_PORT=8080"
Environment="POLL_INTERVAL=5"
ExecStart=/usr/bin/python3 /opt/services/client.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
