# solo-container-and-classics-docker/approach1-systemd/server.service
[Unit]
Description=HTTP Server Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/services
ExecStart=/usr/bin/python3 /opt/services/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
